# zirconium/net-misc/openssh/files/sshdat.service
[Unit]
Description=OpenSSH Per-Connection Server Daemon

[Service]
ExecStart=-/usr/sbin/sshd -i -e
StandardInput=socket
StandardError=journal
